# deploy/npmguard-webhook.service
[Unit]
Description=NpmGuard GitHub Webhook Listener
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/NpmGuard/deploy
ExecStart=/usr/bin/node /root/NpmGuard/deploy/webhook-server.js
Environment=GITHUB_WEBHOOK_SECRET=REPLACE_ME
Restart=always
RestartSec=10
StandardOutput=append:/var/log/npmguard-webhook.log
StandardError=append:/var/log/npmguard-webhook.log
SyslogIdentifier=npmguard-webhook

[Install]
WantedBy=multi-user.target
