# deploy/npmguard-dev.service
[Unit]
Description=NpmGuard Python audit engine (development)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/root/NpmGuard-dev/engine
EnvironmentFile=/root/NpmGuard-dev/engine/.env
ExecStartPre=/root/.local/bin/uv run --frozen alembic upgrade head
ExecStart=/root/.local/bin/uv run --frozen uvicorn npmguard.api:app --host 0.0.0.0 --port 8100
Restart=on-failure
RestartSec=5
TimeoutStopSec=90

[Install]
WantedBy=multi-user.target
